# deb/DEBIAN/postinst
#!/bin/sh

. /usr/share/debconf/confmodule

set -e

pip3 install --system boto3 jmespath

echo ";Generated from debconf" > /etc/lkprint.conf

echo "[Printer]" >> /etc/lkprint.conf
db_get lkprint/printerport
echo "Port = $RET" >> /etc/lkprint.conf

echo "[API]" >> /etc/lkprint.conf
db_get lkprint/apiurl
echo "URL = $RET" >> /etc/lkprint.conf
db_get lkprint/apitoken
echo "Token = $RET" >> /etc/lkprint.conf

echo "[default]" >> /opt/lkprint/.aws/credentials
db_get lkprint/awsid
echo "aws_access_key_id=$RET" >> /opt/lkprint/.aws/credentials
db_get lkprint/awssecret
echo "aws_secret_access_key=$RET" >> /opt/lkprint/.aws/credentials
chown -R root:lkprint /opt/lkprint
chmod 750 /opt/lkprint/.aws/credentials

adduser --system --home /opt/lkprint lkprint
getent group lkprint || groupadd lkprint
usermod -a -G lkprint lkprint
usermod -a -G dialout lkprint
systemctl daemon-reload
systemctl enable lkprint
systemctl start lkprint
