# deb/DEBIAN/prerm
#!/bin/bash
set -e
systemctl disable lkprint
systemctl stop lkprint
deluser lkprint
